# prime-select.service
[Unit]
Description=Prime Select Service
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=prime-select apply-current

[Install]
WantedBy=graphical.target
