# prime-boot-selector.service
[Unit]
Description=Prime Select Service (Boot Selector)
After=multi-user.target
Before=graphical.target

[Service]
Type=oneshot
ExecStart=prime-select prime_booting

[Install]
WantedBy=multi-user.target
