# prime-select.service
[Unit]
Description=Prime Select Service
After=multi-user.target
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=prime-select apply-current

[Install]
WantedBy=multi-user.target
